# .crontab
00 * * * * . $HOME/.cron_functions.sh; notify_time
