# .crontab
00 * * * * . $HOME/.cron_functions.sh; notify_time

* * * * * . $HOME/.cron_functions.sh; run_xmodmap $HOME/.Xmodmap
